# quagga/supervisord.conf
[supervisord]
nodaemon=true

[program:zebra]
command=/usr/local/sbin/zebra -f /usr/local/etc/zebra.conf -i /var/run/quagga/zebra.pid
